# overlay/rc.local
#! /bin/bash
systemctl start djangos
exit 0
